# flex/mcp_server.service
# Flex MCP Server - systemd user service
#
# NOTE: `flex init` generates this unit dynamically with the correct
# Python path. This file is a reference template only.
#
# Install via flex init (recommended):
#   flex init
#
# Manual install:
#   cp flex/mcp_server.service ~/.config/systemd/user/flex-mcp.service
#   systemctl --user daemon-reload
#   systemctl --user enable --now flex-mcp
#
# View logs:
#   journalctl --user -u flex-mcp -f

[Unit]
Description=Flex MCP Server
After=network.target

[Service]
Type=simple
ExecStart=python -m flex.mcp_server --http --port 7134
Restart=always
RestartSec=5
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=default.target
